# src/ign-api/crontab
5 */6    * * *   root    . /usr/bin/bash /app/run-worker.sh
